# scripts-and-test/vps-worker/juleha-worker.service
[Unit]
Description=Juleha Enrichment Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/aicenghub
EnvironmentFile=/etc/aicenghub/worker.env
ExecStart=/usr/bin/node scripts-and-test/vps-worker/worker.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

